# deploy/systemd/tgn-watch.service
[Unit]
Description=tgn-watch lightweight monitoring agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tgn-watch
Group=tgn-watch
ExecStart=/usr/local/bin/tgn-watch -config /etc/tgn-watch/config.yml
Restart=on-failure
RestartSec=3s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/tgn-watch

[Install]
WantedBy=multi-user.target
